# backend/systemd/smartcam.service
[Unit]
Description=SmartCam Soccer Main Service
After=network.target
Wants=smartcam-manager.service smartcam-status.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartCam-Soccer/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/SmartCam-Soccer/backend/venv/bin/python main.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=smartcam

[Install]
WantedBy=multi-user.target 