# backend/systemd/smartcam-manager.service
[Unit]
Description=SmartCam Soccer Manager Service
After=network.target
PartOf=smartcam.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartCam-Soccer/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/SmartCam-Soccer/backend/venv/bin/python manage_cameras.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=smartcam-manager

[Install]
WantedBy=multi-user.target 