# backend/systemd/smartcam-status.service
[Unit]
Description=SmartCam Soccer Status Service
After=network.target
PartOf=smartcam.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/SmartCam-Soccer/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/SmartCam-Soccer/backend/venv/bin/python update_system_status.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=smartcam-status

[Install]
WantedBy=multi-user.target 